# scripts/cockpit_agent_git/cockpit-agent-git.service
[Unit]
Description=Git Agent for Cockpit
After=network.target

[Service]
Type=simple
User=cockpit-agent
Group=docker
WorkingDirectory=/opt/cockpit-agent-git
ExecStart=/opt/cockpit-agent-git/venv/bin/python agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
EnvironmentFile=/opt/cockpit-agent-git/.env

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
